# etc/supervisord.d/mongorestore.ini
[program:mongorestore]
command=/opt/qnib/mongodb/bin/restore.sh
startsecs=0 
redirect_stderr=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
