# etc/supervisord.d/mongodb.ini
[program:mongodb]
command=/opt/qnib/mongodb/bin/start.sh
redirect_stderr=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
